# systemd/fixtures/manager_dummy.service
[Unit]
Description=dummy unit for Mirasol e2e tests

[Service]
ExecStart=/bin/sleep 400
